# rootfs/etc/service/geolite2/run
#!/bin/sh -e
exec 2>&1
(cd /opt/app; exec ./geolite2)